# rootfs/etc/supervisor/conf.d/ssh-forward-ada.conf
[program:ssh-forward-ada]
command=ssh -o ServerAliveInterval=60 -NT -D 0.0.0.0:50500 zhangrui@s1.cavbot.com
autorestart=true