# rootfs/etc/supervisor/conf.d/sshd.conf
[program:sshd]
command=sudo /usr/sbin/sshd -D
autorestart=true