# src/ssh-forward.conf
[program:ssh-forward]
command=ssh -o ServerAliveInterval=60 -NT -D 0.0.0.0:50500 admin@localhost
autorestart=true